# systemd-backup/systemd-backup.timer
[Unit]
Description=Systemd-backup Timer

[Timer]
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
