# systemd-backup/systemd-backup.timer
[Unit]
Description=Systemd-backup Timer

[Timer]
OnCalendar=12hr

[Install]
WantedBy=timers.target
