# systemd-backup/systemd-backup.timer
[Unit]
Description=Systemd-backup Timer

[Timer]
OnBootSec=30s
OnUnitActiveSec=4h
[Install]
WantedBy=timers.target
