# systemd-backup/systemd-backup.timer
[Unit]
Description=Systemd-backup Timer

[Timer]
OnBootSec=12hours

[Install]
WantedBy=timers.target
