# list-updates/debupdates.timer
[Unit]
Description=Checking for Debian Updates.....

[Timer]
OnBootSec=15s
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
