# update-service/debupdates.timer
[Unit]
Description=Checking for Debian Updates.....

[Timer]
OnBootSec=15s
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
