# mount_media/home-tknomanzr-remote.automount
[Automount]
Where=/home/tknomanzr/remote
DirectoryMode=0775

[Install]
WantedBy=multi-user.target
