# snapshot-cleanup/snapshot-cleanup.timer
[Unit]
Description=Daily cleanup of apt-btrfs-snapper snapshots

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

